# deploy/services/gddo-server/run
#!/bin/sh
exec /gopath/bin/gddo-server 2>&1
